# new/resolvconf.service
[Unit]
Before=network-pre.target
Description=resolvconf update
Wants=network-pre.target
X-Restart-Triggers=/nix/store/v1babmp3lkcba5733i4rq99mk8ha2jd0-X-Restart-Triggers-resolvconf

[Service]
Environment="LOCALE_ARCHIVE=/nix/store/lw4a8r3lm5irklghq3pjj9spj9bvzwkd-glibc-locales-2.40-218/lib/locale/locale-archive"
Environment="PATH=/nix/store/wkkwxc04gdw6b263l1h29pjarjnjdyb6-coreutils-9.8/bin:/nix/store/bb1ipjfz8g6cir41rr0zc2x2ah43kzar-findutils-4.10.0/bin:/nix/store/d4c56s8wa6rz2dnw6ridg7r1dvax1gky-gnugrep-3.12/bin:/nix/store/0zqnw11n1hk8mflzzvrz7sv1rm1cbnp8-gnused-4.9/bin:/nix/store/bigkpra9jw48fip69q4wndsf4kb3d2w9-systemd-258.3/bin:/nix/store/wkkwxc04gdw6b263l1h29pjarjnjdyb6-coreutils-9.8/sbin:/nix/store/bb1ipjfz8g6cir41rr0zc2x2ah43kzar-findutils-4.10.0/sbin:/nix/store/d4c56s8wa6rz2dnw6ridg7r1dvax1gky-gnugrep-3.12/sbin:/nix/store/0zqnw11n1hk8mflzzvrz7sv1rm1cbnp8-gnused-4.9/sbin:/nix/store/bigkpra9jw48fip69q4wndsf4kb3d2w9-systemd-258.3/sbin"
Environment="TZDIR=/nix/store/8i1vdm6jfvvchysfbymf1gxd54danazj-tzdata-2025c/share/zoneinfo"
ExecStart=/nix/store/mh2b8i437pk76wpzcds38313w0i0hpkm-unit-script-resolvconf-start/bin/resolvconf-start 
RemainAfterExit=true
Type=oneshot

[Install]
WantedBy=multi-user.target
